# deploy/systemd/hacs-stats-web.service
[Unit]
Description=hacs-stats web server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=hacs-stats
Group=hacs-stats
WorkingDirectory=/opt/hacs-stats/apps/web
EnvironmentFile=/etc/hacs-stats/env
# Invoke node directly — pnpm 11 runs a "verify deps before run" step that
# tries to write under WorkingDirectory, which collides with the
# ProtectSystem=strict hardening below. node + tsx loader does the same
# work without the package-manager scaffolding.
ExecStart=/usr/bin/node --import tsx src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening — Node has no business writing outside /var/lib/hacs-stats or
# reading /etc beyond its own envfile.
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hacs-stats
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false

[Install]
WantedBy=multi-user.target
