# deploy/systemd/hacs-stats-scrape.service
[Unit]
Description=hacs-stats — one-shot scrape job
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=hacs-stats
Group=hacs-stats
WorkingDirectory=/opt/hacs-stats/apps/scraper
EnvironmentFile=/etc/hacs-stats/env
# See hacs-stats-web.service for why node-directly instead of pnpm.
ExecStart=/usr/bin/node --import tsx src/index.ts
StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hacs-stats
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
