# deploy/systemd/hacs-stats-discover.service
[Unit]
Description=hacs-stats — weekly discovery sweep (15-band GitHub code-search for hacs.json)
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=hacs-stats
Group=hacs-stats
WorkingDirectory=/opt/hacs-stats
EnvironmentFile=/etc/hacs-stats/env
# scripts/discover-bands.ts walks 15 size-banded queries in one process
# to break past GitHub's 1000-result code-search cap. The earlier
# single-query version (scripts/discover.ts) hit the cap on the busy
# tail of the catalogue and missed most candidates after the first run.
# Bands take 25-40 min wall-clock; well within the weekly cadence.
ExecStart=/usr/bin/node --import tsx scripts/discover-bands.ts
TimeoutStartSec=2700
StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/hacs-stats
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
